# config_files/crontab

15,45  *  *   *   *     /home/trakkor/rails_app/fetch.sh
49-59 0,12,20  *   *   *     /home/trakkor/rails_app/backup-database.sh

